# deploy/crontab
# Schedule for the scheduler service. Times are UTC (container clock).
# SHELL+BASH_ENV make each job a bash shell that restores the container env
# (cron strips it); >/proc/1/fd/{1,2} sends output to `docker logs`.
SHELL=/bin/bash
BASH_ENV=/tmp/container-env.sh

# Each line maps a TIME to a MANIFEST path. The active line runs a shipped sample
# (examples/, baked in) so the scheduler works out of the box. For your own
# manifests: build/test one with `paramify tui` (it writes to manifests/), rebuild
# the image so it's baked in, then point a line at manifests/<name>.yaml.
#
# min  hour  dom  mon  dow   command
0     2     *    *    *      cd /app && ./deploy/run-and-upload.sh examples/minimal_run.yaml >/proc/1/fd/1 2>/proc/1/fd/2
# 0   3     *    *    1      cd /app && ./deploy/run-and-upload.sh manifests/your-manifest.yaml >/proc/1/fd/1 2>/proc/1/fd/2
